# SIPp/server/supervisord.conf
[supervisord]
nodaemon=true

[program:sipp-server]
command=sh /opt/server_start.sh