# SIPp/client/supervisord.conf
[supervisord]
nodaemon=true

[program:sipp-client]
command=sh /opt/client_start.sh
